# deploy/marketingagents.service
[Unit]
Description=MarketingAgents Streamlit Backend
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/MarketingAgents/backend
ExecStart=/home/ubuntu/MarketingAgents/backend/.venv/bin/streamlit run /home/ubuntu/MarketingAgents/app.py --server.port=8502 --server.address=127.0.0.1
Restart=always
RestartSec=5
Environment="PATH=/home/ubuntu/MarketingAgents/backend/.venv/bin"
Environment="PYTHONPATH=/home/ubuntu/MarketingAgents/backend"

[Install]
WantedBy=multi-user.target
